# systemd/hwt905-raspi-scheduled.service
[Unit]
Description=HWT905 Raspberry Pi IMU Data Processing Service (Scheduled Mode)
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=trandung
Group=trandung
WorkingDirectory=/home/trandung/Data/Aitogy Projects/hwt905-raspi
Environment=PYTHONPATH=/home/trandung/Data/Aitogy Projects/hwt905-raspi
ExecStart=/usr/bin/python3 /home/trandung/Data/Aitogy Projects/hwt905-raspi/main.py --mode scheduled
StandardOutput=journal
StandardError=journal

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
