# systemd/hwt905-raspi.service
[Unit]
Description=HWT905 Raspberry Pi IMU Data Processing Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=trandung
Group=trandung
WorkingDirectory=/home/trandung/Data/Aitogy Projects/hwt905-raspi
Environment=PYTHONPATH=/home/trandung/Data/Aitogy Projects/hwt905-raspi
# Mặc định chạy ở chế độ realtime
# Để thay đổi chế độ, sửa dòng ExecStart bên dưới:
# --mode realtime: gửi dữ liệu ngay lập tức
# --mode batch: gửi dữ liệu theo batch
# --mode scheduled: gửi dữ liệu theo lịch trình
ExecStart=/usr/bin/python3 /home/trandung/Data/Aitogy Projects/hwt905-raspi/main.py --mode realtime
StandardOutput=journal
StandardError=journal

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
